# systemd/publish-images.path
[Unit]
Description=Watch for published radar images

[Path]
PathChanged=/media/holding/nws/nws-publish
Unit=publish-images.service

[Install]
WantedBy=multi-user.target